# infra/provisioning/files/heub.conf
start on startup
chdir /opt/heub/
mkdir -p /opt/heub/log
exec /opt/heub/bin/heub -p 8081 -e prod +RTS -A4M -qg1
